# doc/examples/etc/cron.daily/snazzer
#!/bin/sh
set -e
. /etc/default/snazzer
snazzer --all
snazzer --prune --force --all
snazzer --measure --all
