# doc/examples/etc/cron.daily/snazzer-receive
#!/bin/sh
set -e
su receiveruser -c "/bin/sh snazzer-receive.sh"
